# deploy/systemd/terminus-discord.service
[Unit]
Description=Terminus-OS Discord tenant (k12-personal)
Documentation=https://github.com/TokenAmby-Code/Terminus-OS
# The tenant reaches Discord over the network; it has no Token-API/tmuxctld dependency and no
# NAS mount dependency (the old daemon's launchd ordering against those is gone).
After=network-online.target
Wants=network-online.target
# Crash-loop backstop: after 5 failures inside 5 minutes, stop and stay failed (loud), rather
# than hot-spinning. A persistently missing credential or the not-yet-installed transport
# adapter lands the unit in `failed`, visible to the operator — never a silent restart churn.
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=tokenamby
# Credentials are ENV-ONLY. Token-Fleet provisions this file on the host from its secret store
# (root:root, 0600); it is NEVER in the checkout. Every variable it must carry is enumerated in
# deploy/env/terminus-discord.env.example. This replaces the old daemon's macOS Keychain reads.
EnvironmentFile=/etc/terminus-os/discord.env
# The tenant runs from the LOCAL runtime checkout — the hot path never traverses the NAS/SMB
# mount (recon §6). NAS stays durable coordination, not a runtime dependency.
WorkingDirectory=/home/tokenamby/runtimes/Terminus-OS/live
ExecStart=/usr/bin/env bun run packages/discord/src/main.ts
Restart=on-failure
RestartSec=10
# Logs go to the journal (journalctl -u terminus-discord). No macOS `say`, no Keychain, no
# Token-API/tmuxctld coupling, no name-based Inquisition self-test.
StandardOutput=journal
StandardError=journal

# Hardening: least privilege for a network client that only reads its own checkout.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
